# ffhb/ffhb-breminale/files/etc/init.d/watchdog
#!/bin/sh /etc/rc.common

START=99
USE_PROCD=1
BIN=/usr/sbin/watchdog
