# ffhb/ffhb-breminale/files/etc/hotplug.d/net/20-alfred
#!/bin/sh

if [ "$ACTION" = "add" -a "$DEVICENAME" = "br-client" ]; then
  /etc/init.d/alfred restart
fi
